# nftables/debian/lib/systemd/system/nftables.service.d/sn-nftables.conf
[Service]
ExecStop=
ExecStop=/bin/sh /usr/share/solarnode/bin/nftables-utils.sh reset

