# solarnode-base-java17/debian/lib/systemd/system/solarnode.service.d/solarnode-base-java17.conf
[Service]
Environment="OSGI_PROFILE_OPTS=-Dosgi.java.profile=file:conf/java17-server.profile"
