# mobile-usb-ppp/debian/lib/systemd/system/sn-mobile-usb-ppp-pppd.service
[Unit]
Description=USB mobile network PPPd managed connection

[Service]
ExecStartPre=-/bin/rm -f /var/lock/LCK..modem
ExecStart=/usr/sbin/pppd call sn-provider nodetach nolog
Restart=on-failure
RestartSec=30
