# solarnode-bluetooth-setup/debian/etc/systemd/system/bluetooth.service.d/solarnode-bluetooth-setup.conf
[Service]
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd -P battery

