# pi-mobile-shield-usb/debian/lib/systemd/system/sn-mobile-shield-usb-reconnect.service
[Unit]
Description=USB mobile network auto reconnect

[Service]
Type=oneshot
ExecStart=/usr/share/solarnode/bin/mobile-shield-usb-reconnect.sh
KillMode=none
