# solarnode-base-log4j2/debian/lib/systemd/system/solarnode.service.d/solarnode-base-log4j2.conf
[Service]
Environment=LOG4J_CONFIGURATION_FILE=/etc/solarnode/log4j2.xml
Environment="SOLARNODE_LOGGING_PROPS=-Dlog4j2.contextSelector=org.apache.logging.log4j.core.selector.BasicContextSelector \
	-Dlog4j2.provider=org.apache.logging.log4j.core.impl.Log4jProvider \
	-Dslf4j.provider=org.apache.logging.slf4j.SLF4JServiceProvider"
