# pi-mobile-shield-quectel/debian/sn-pi-mobile-shield-quectel.prerm
#!/bin/sh -e

service_down () {
	name="$1"
	if systemctl is-enabled "$name" >/dev/null; then
		echo "Disabling $name..."
		systemctl disable --now "$name" || true
	fi
}

if [ "$1" = "remove" ]; then
	sn-mobile-shield-quectel-reconnect.timer
	service_down sn-mobile-shield-quectel-pppd
fi
