# solarnode-base-java11/debian/lib/systemd/system/solarnode.service.d/solarnode-base-java11.conf
[Service]
Environment="OSGI_PROFILE_OPTS=-Dosgi.java.profile=file:conf/java11-server.profile"
