# pi-bluetooth/debian/etc/systemd/system/bluetooth.service.d/sn-bluetooth.conf
[Service]
ExecStart=
ExecStart=/usr/lib/bluetooth/bluetoothd -C --noplugin=sap
ExecStartPost=/usr/bin/sdptool add SP