# solarnode-bluetooth-setup/debian/solarnode-bluetooth-setup.postinst
#!/bin/sh
set -e

CONF="/usr/share/solarnode/default/solarnode-bluetooth-setup"
VENDOR_CONF="/etc/default/solarnode-bluetooth-setup"

if [ -e "$CONF" ]; then
	. "$CONF"
fi
if [ -e "$VENDOR_CONF" ]; then
	. "$VENDOR_CONF"
fi

service_up () {
	name="$1"
	start="$2"
	if ! systemctl is-enabled "$name" >/dev/null; then
		echo "Enabling $name..."
		systemctl enable "$name" || true
	fi
	if [ -n "$start" ]; then
		if ! systemctl is-active "$name" >/dev/null; then
			systemctl start "$name" || true
		else
			systemctl restart "$name" || true
		fi
	fi
}

BLUETOOTH_CONF="/etc/bluetooth/main.conf"
if grep -q '^DiscoverableTimeout = 0' $BLUETOOTH_CONF; then
	true
else
	echo "Setting DiscoverableTimeout = 0 in $BLUETOOTH_CONF."
	if grep -q 'DiscoverableTimeout[[:space:]]*=' $BLUETOOTH_CONF; then
		# update
		sed -i -e '/DiscoverableTimeout[[:space:]]*=/c DiscoverableTimeout = 0' $BLUETOOTH_CONF || true
	else
		# add in
		echo 'DiscoverableTimeout = 0' >>$BLUETOOTH_CONF
	fi
fi
if grep -q '^AutoEnable[[:space:]]*=[[:space:]]*true' $BLUETOOTH_CONF; then
	true
else
	echo "Setting AutoEnable = true in $BLUETOOTH_CONF."
	if grep -q 'AutoEnable[[:space:]]*=' $BLUETOOTH_CONF; then
		# update
		sed -i -e '/AutoEnable[[:space:]]*=/c AutoEnable = true' $BLUETOOTH_CONF || true
	else
		# add in
		echo 'AutoEnable = true' >>$BLUETOOTH_CONF
	fi
fi

setup_user_conf_example='/usr/share/solarnode/example/solarnode-bluetooth-setup-user.csv'
setup_user_conf='/etc/solarnode/auto-settings.d/solarnode-bluetooth-setup.csv'
if [ -z "$CFG_WITHOUT_SETUP_USER" -a ! -e "${setup_user_conf}" ]; then
	if [ -e "${setup_user_conf_example}" ]; then
		echo "Creating default setup user configuration ${setup_user_conf}"
		cp "${setup_user_conf_example}" "${setup_user_conf}"
	fi
fi

systemctl daemon-reload || true
service_up solarnode-bt-setup-peripheral true
