# solarnode-base-console/debian/lib/systemd/system/solarnode.service.d/solarnode-base-console.conf
[Service]
Environment="OSGI_CONSOLE_OPTS=-console 4202"
