# pi-mobile-shield-quectel/debian/lib/systemd/system/sn-mobile-shield-quectel-reconnect.service
[Unit]
Description=Quectel mobile network auto reconnect

[Service]
Type=oneshot
ExecStart=/usr/share/solarnode/bin/mobile-shield-quectel-reconnect.sh
KillMode=none
