# mobile-mm/debian/lib/systemd/system/sn-mobile-mm-init.path
[Unit]
Description=Mobile ModemManager setup modem device

[Path]
Unit=sn-mobile-mm-init.service
PathExists=/dev/modem

[Install]
WantedBy=multi-user.target
