# pypowerwall/debian/lib/systemd/system/sn-pypowerwall-proxy.path
[Unit]
Description=SolarNode PyPowerwall Proxy configuration
Wants=network.target

[Path]
Unit=sn-pypowerwall-proxy.service
PathChanged=/etc/solarnode/sn-pypowerwall-proxy.conf

[Install]
Alias=multi-user.target.wants/sn-pypowerwall-proxy.path
