# solarnode-app-camera-ffmpeg/debian/lib/systemd/system/solarnode-app-camera-ffmpeg-cleaner.timer
[Unit]
Description=Clean old FFmpeg camera snapshot images daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
