# mbpoll-reset/debian/lib/systemd/system/sn-mbpoll-reset@.timer
[Unit]
Description=Reset Modbus hourly

[Timer]
OnCalendar=hourly
AccuracySec=100ms
RandomizedDelaySec=0

[Install]
WantedBy=timers.target
