# heartbeat-log/debian/lib/systemd/system/sn-heartbeat-log.timer
[Unit]
Description=Generate heartbeat messages

[Timer]
OnBootSec=10s
OnUnitActiveSec=5min
AccuracySec=2s

[Install]
WantedBy=timers.target
